# 乐学偶得资料/william_web_server.ini
[program:william_web_server]
directory = /root/william_flask_edu/william_web_server
command = /root/william_flask_edu/william_web_server/venv/bin/gunicorn -w 5 run:app
user = root
autostart = true
autorestart = true
stopasgroup = true
killasgroup = true
stderr_logfile = /var/log/william_web_server/william_web_server.err.log
stdout_logfile = /var/log/william_web_server/william_web_server.out.log
